# contrib/aur/z13ctl-bin.install
post_install() {
    z13ctl setup
    systemctl --global enable z13ctl.socket z13ctl.service || true
}

post_upgrade() {
    udevadm control --reload-rules
    udevadm trigger
    systemctl daemon-reload || true
}
