# contrib/aur/z13ctl-bin.install
post_install() {
    udevadm control --reload-rules
    udevadm trigger
    systemctl enable --now z13ctl-perms.service
    systemctl --global enable z13ctl.socket z13ctl.service
}

post_upgrade() {
    udevadm control --reload-rules
    udevadm trigger
}
